# supervisord-virtuoso.conf
[program:virtuoso]
command=/start-virtuoso.sh
numprocs=1
autostart=true
autorestart=true