# supervisord-tomcat6.conf
[program:tomcat6]
command=/start-tomcat6.sh
numprocs=1
autostart=true
autorestart=true